# pkgs/mock-monitor/DEBIAN/prerm
#!/bin/bash
# Pre-removal script for mock-monitor

set -e

# Stop the service if it's running
if systemctl is-active --quiet mock-monitor.service; then
    echo "Stopping mock-monitor service..."
    systemctl stop mock-monitor.service
fi

# Disable the service
if systemctl is-enabled --quiet mock-monitor.service; then
    echo "Disabling mock-monitor service..."
    systemctl disable mock-monitor.service
fi