# pkgs/zfs-helper-client/usr/share/doc/zfs-helper-client/examples/backup@.service
[Unit]
Description=Snapshot before job %i

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'zfs-helperctl snapshot tank/home/${USER}@%i'
